# scheduler/com.insightque.stock-alpaca.analyst-review.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.insightque.stock-alpaca.analyst-review</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/insightque/stock-alpaca/scripts/run-analyst-review-codex.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/insightque/stock-alpaca</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>20</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/insightque/stock-alpaca/logs/analyst-review.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/insightque/stock-alpaca/logs/analyst-review.err.log</string>
</dict>
</plist>
