# scheduler/com.insightque.stock-alpaca.hourly-autopilot.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.insightque.stock-alpaca.hourly-autopilot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/insightque/stock-alpaca/scripts/run-hourly-autopilot-codex.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/insightque/stock-alpaca</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/Library/Frameworks/Python.framework/Versions/3.11/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/insightque</string>
    <key>LANG</key>
    <string>en_US.UTF-8</string>
    <key>LC_ALL</key>
    <string>en_US.UTF-8</string>
    <key>NO_PROXY</key>
    <string>127.0.0.1,localhost</string>
    <key>CODEX_SCHEDULED_CODEX_HOME</key>
    <string>/Users/insightque/.codex</string>
    <key>CODEX_AUTOPILOT_RESEARCH_CACHE_DIR</key>
    <string>/Users/insightque/stock-alpaca/.cache/research-mcp-preflight</string>
    <key>CODEX_AUTOPILOT_RESEARCH_CIRCUIT_SECONDS</key>
    <string>600</string>
    <key>CODEX_AUTOPILOT_REGISTER_RESEARCH_MCP</key>
    <string>auto</string>
  </dict>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>22</integer>
      <key>Minute</key>
      <integer>31</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>22</integer>
      <key>Minute</key>
      <integer>51</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>23</integer>
      <key>Minute</key>
      <integer>11</integer>
    </dict>
    <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>31</integer></dict>
    <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>51</integer></dict>
    <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>11</integer></dict>
    <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>31</integer></dict>
    <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>51</integer></dict>
    <dict><key>Hour</key><integer>1</integer><key>Minute</key><integer>11</integer></dict>
    <dict><key>Hour</key><integer>1</integer><key>Minute</key><integer>31</integer></dict>
    <dict><key>Hour</key><integer>1</integer><key>Minute</key><integer>51</integer></dict>
    <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>11</integer></dict>
    <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>31</integer></dict>
    <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>51</integer></dict>
    <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>11</integer></dict>
    <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>31</integer></dict>
    <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>51</integer></dict>
    <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>11</integer></dict>
    <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>31</integer></dict>
    <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>51</integer></dict>
    <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>11</integer></dict>
    <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>31</integer></dict>
    <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>51</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/insightque/stock-alpaca/logs/hourly-autopilot.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/insightque/stock-alpaca/logs/hourly-autopilot.err.log</string>
</dict>
</plist>
